# client/client/rc.local
# Put this in the rc.local file:
#   sudo nano /etc/rc.local

sudo /usr/bin/ruby /home/pi/rooster/client.rb &